# deploy/qmdsr.service
[Unit]
Description=qmdsr - qmd search router middleware
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/qmdsr -config /etc/qmdsr/qmdsr.yaml
Restart=always
RestartSec=5
OOMPolicy=continue
Environment=HOME=/var/lib/qmdsr
WorkingDirectory=/var/lib/qmdsr
StateDirectory=qmdsr
RuntimeDirectory=qmdsr
LogsDirectory=qmdsr
WatchdogSec=120

[Install]
WantedBy=multi-user.target
